# systemd/infra-demo.service
[Unit]
Description=Infrastructure Demo HTTP Health Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/infra-demo
ExecStart=/usr/bin/python3 /opt/infra-demo/service.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
